# add_files/rk3568_fs/etc/udev/rules.d/99-set-mac.rules
ACTION=="add", SUBSYSTEM=="net", KERNEL=="e*", RUN+="/usr/local/bin/set_mac.sh %k"
